# .config/systemd/user.control/clash.service
[Unit]
Description=clash daemon

[Service]
ExecStart=/home/wangjie/Gopath/bin/clash
Restart=on-failure

[Install]
WantedBy=default.target

