# .config/systemd/user.control/xcompmgr.service
[Unit]
Description=xcompmgr daemon

[Service]
ExecStart=xcompmgr
Restart=on-failure

[Install]
WantedBy=default.target

