# .config/systemd/user.control/xscreensaver.service
[Unit]
Description=xscreensaver daemon

[Service]
ExecStart=/usr/bin/xscreensaver -no-splash
Restart=on-failure

[Install]
WantedBy=default.target

